# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for DOMAIN

[Service]
Restart=on-failure
User=USER
WorkingDirectory=/home/USER/sites/DOMAIN
EnvironmentFile=/home/USER/sites/DOMAIN/.env

ExecStart=/home/USER/sites/DOMAIN/virtualenv/bin/gunicorn \
    --bind unix:/tmp/DOMAIN.socket \
    superlists.wsgi

[Install]
WantedBy=multi-user.target
